# etc/supervisord.d/craft-worker.ini
[program:craft-queue]
command=/usr/bin/nice -n 10 /usr/bin/php /var/www/html/craft queue/listen
autostart=true
autorestart=true

[program:craft-schedules]
command=/usr/bin/nice -n 10 /usr/bin/php /var/www/html/craft schedules/listen
autostart=true
autorestart=true