# etc/supervisord.d/nginx.ini
[program:nginx]
command=sudo /usr/sbin/nginx -c /etc/nginx/nginx.conf
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=0
autorestart=true
