# etc/supervisord.d/php-fpm.ini
[program:php-fpm]
command=sudo /usr/sbin/php-fpm -c /etc/php-fpm.conf
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=0
autorestart=true